# cloud/roles/gpumon/templates/gpumon.service
[Unit]
Description=Collect GPU metrics

[Service]
User=nobody
Group=nogroup
Type=oneshot
ExecStart=/usr/local/bin/gpumon
